# mcp/cobbleai-mcp.service
[Unit]
Description=CobbleAI MCP Tool Server
After=network.target

[Service]
Type=simple
User=administrator
WorkingDirectory=/home/administrator/cobbleai/mcp
EnvironmentFile=/home/administrator/cobbleai/.env
ExecStart=/home/administrator/cobbleai/venv/bin/python server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
